# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisor/supervisord.pid
childlogdir=/var/log/supervisor

[program:comfyui]
command=/usr/bin/python3.12 /opt/comfyui/main.py --listen 0.0.0.0 --port 8188 --cuda-device 1 --preview-method auto
directory=/opt/comfyui
environment=DISPLAY=":1",CUDA_VISIBLE_DEVICES="1"
autorestart=true
priority=40
startretries=5
startsecs=10
stdout_logfile=/var/log/supervisor/comfyui.log
stderr_logfile=/var/log/supervisor/comfyui-err.log

[program:web-interface]
command=/usr/bin/python3.12 /opt/web-interface/app.py
directory=/opt/web-interface
autorestart=true
priority=50
stdout_logfile=/var/log/supervisor/web-interface.log
stderr_logfile=/var/log/supervisor/web-interface-err.log

[program:terminal]
command=/usr/local/bin/ttyd --port 7681 --writable bash
directory=/opt/gaussian-toolkit
autorestart=true
priority=60
startretries=3
startsecs=2
stdout_logfile=/var/log/supervisor/terminal.log
stderr_logfile=/var/log/supervisor/terminal-err.log
